# crontab
0 * * * * build-index
30 0 * * * build-sitemap
15 */6 * * * load-news-releases
5 */6 * * * load-blueprints
0 */6 * * * load-wikis
